# grocery_ai_planning_ws/devel/etc/catkin/profile.d/cmake.lock
/home/rodri/grocery_ai_planning_ws/devel/.private/grocery/etc/catkin/profile.d/cmake.lock